# script.cron
* * * * * /usr/bin/sh /home/ubuntu/editDate.sh
* * * * * sleep 30; /usr/bin/sh /home/ubuntu/editDate.sh
